# sonarr.conf
[program:sonarr]
user = nobody
command = /usr/sbin/sonarr